# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid

[program:vnc]
command=/opt/scripts/start-vnc.sh
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/vnc.log
stderr_logfile=/var/log/vnc.err.log
startsecs=2

[program:wechat]
command=/opt/scripts/start-wechat.sh
autostart=true
autorestart=unexpected
priority=20
stdout_logfile=/var/log/wechat.log
stderr_logfile=/var/log/wechat.err.log
startsecs=5
environment=DISPLAY=":1",WINEPREFIX="%(ENV_WINEPREFIX)s"

[program:wcf-bridge]
command=/opt/scripts/start-wcf.sh
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/wcf-bridge.log
stderr_logfile=/var/log/wcf-bridge.err.log
startsecs=10
